# eval/supervisor/muon_glue_completion.conf
[program:muon-glue-boolq]
command=/workspace/nanoGPT/eval/run_muon_glue_completion.sh boolq
directory=/workspace/nanoGPT
environment=CUDA_VISIBLE_DEVICES="0",PYTORCH_CUDA_ALLOC_CONF="expandable_segments:True",PROC_NAME="%(program_name)s"
autostart=false
autorestart=false
startsecs=1
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:muon-glue-wsc]
command=/workspace/nanoGPT/eval/run_muon_glue_completion.sh wsc
directory=/workspace/nanoGPT
environment=CUDA_VISIBLE_DEVICES="1",PYTORCH_CUDA_ALLOC_CONF="expandable_segments:True",PROC_NAME="%(program_name)s"
autostart=false
autorestart=false
startsecs=1
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:muon-glue-mrpc]
command=/workspace/nanoGPT/eval/run_muon_glue_completion.sh mrpc
directory=/workspace/nanoGPT
environment=CUDA_VISIBLE_DEVICES="2",PYTORCH_CUDA_ALLOC_CONF="expandable_segments:True",PROC_NAME="%(program_name)s"
autostart=false
autorestart=false
startsecs=1
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:muon-glue-multirc]
command=/workspace/nanoGPT/eval/run_muon_glue_completion.sh multirc
directory=/workspace/nanoGPT
environment=CUDA_VISIBLE_DEVICES="3",PYTORCH_CUDA_ALLOC_CONF="expandable_segments:True",PROC_NAME="%(program_name)s"
autostart=false
autorestart=false
startsecs=1
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[group:muon-glue-completion]
programs=muon-glue-boolq,muon-glue-wsc,muon-glue-mrpc,muon-glue-multirc
priority=510
