# eval/supervisor/muon_lr_pilot.conf
[program:muon-lr-1e-5]
command=/workspace/nanoGPT/eval/run_muon_lr_pilot.sh 1e-5 lr1e-5-multirc-ga2
directory=/workspace/nanoGPT
environment=CUDA_VISIBLE_DEVICES="0",TASK_SET="multirc",MULTIRC_GA="2",PYTORCH_CUDA_ALLOC_CONF="expandable_segments:True",PROC_NAME="%(program_name)s"
autostart=false
autorestart=false
startsecs=1
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:muon-lr-3e-5]
command=/workspace/nanoGPT/eval/run_muon_lr_pilot.sh 3e-5 lr3e-5-multirc-ga2
directory=/workspace/nanoGPT
environment=CUDA_VISIBLE_DEVICES="1",TASK_SET="multirc",MULTIRC_GA="2",PYTORCH_CUDA_ALLOC_CONF="expandable_segments:True",PROC_NAME="%(program_name)s"
autostart=false
autorestart=false
startsecs=1
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:muon-lr-6e-5]
command=/workspace/nanoGPT/eval/run_muon_lr_pilot.sh 6e-5 lr6e-5-multirc-ga2
directory=/workspace/nanoGPT
environment=CUDA_VISIBLE_DEVICES="2",TASK_SET="multirc",MULTIRC_GA="2",PYTORCH_CUDA_ALLOC_CONF="expandable_segments:True",PROC_NAME="%(program_name)s"
autostart=false
autorestart=false
startsecs=1
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:muon-lr-1e-4]
command=/workspace/nanoGPT/eval/run_muon_lr_pilot.sh 1e-4 lr1e-4-multirc-ga2
directory=/workspace/nanoGPT
environment=CUDA_VISIBLE_DEVICES="3",TASK_SET="multirc",MULTIRC_GA="2",PYTORCH_CUDA_ALLOC_CONF="expandable_segments:True",PROC_NAME="%(program_name)s"
autostart=false
autorestart=false
startsecs=1
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[group:muon-lr-pilot]
programs=muon-lr-1e-5,muon-lr-3e-5,muon-lr-6e-5,muon-lr-1e-4
priority=500
