# deploy/supervisor/offdev100m-eval-seed1338-baseline.conf
[program:offdev100m_eval_seed1338_baseline]
directory=/workspace/repo/nanoGPT
command=/bin/bash /workspace/repo/nanoGPT/eval/run_offdev_100m_zero_shot_single.sh
environment=VARIANT="bl100m-d512L32-do0.1-gate-offdev-aoaw28-aoat31-u57-b16ga32-s1338",GPU="2",PY="/workspace/envs/babylm-eval/bin/python",WAIT_FOR_FINAL="0"
autostart=false
autorestart=false
startsecs=5
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=120
redirect_stderr=true
stdout_logfile=/workspace/repo/nanoGPT/logs/eval-offdev-100m/seed1338-baseline.supervisor.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
