# deploy/supervisor/offdev100m-eval-seed1337-attnres.conf
[program:offdev100m_eval_seed1337_attnres]
directory=/workspace/repo/nanoGPT
command=/bin/bash /workspace/repo/nanoGPT/eval/run_offdev_100m_zero_shot_single.sh
environment=VARIANT="bl100m-d512L32-do0.1-gate-attnres8-offdev-aoaw28-aoat31-u57-b16ga32",GPU="1",PY="/workspace/envs/babylm-eval/bin/python",WAIT_FOR_FINAL="1",POLL_SECONDS="30"
autostart=false
autorestart=false
startsecs=5
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=120
redirect_stderr=true
stdout_logfile=/workspace/repo/nanoGPT/logs/eval-offdev-100m/seed1337-attnres.supervisor.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
